# src/deb/control/postinst
#!/bin/sh

echo "\033[36m HI I'M A POSTINST SCRIPT `date +"%s"` \033[39m"
sudo chkconfig quarkus-microservice on
exit 0
